# c720-hdmi-switcher/lib/udev/rules.d/hdmi_sound.rules
SUBSYSTEM=="drm", ACTION=="change", HOTPLUG=="1", RUN+="/usr/bin/hdmi_sound_toggle_all"

